# systemd/midscroll.service
[Unit]
Description=midscroll - Windows-style middle-drag autoscroll

[Service]
ExecStartPre=-/usr/sbin/modprobe uinput
ExecStart=/usr/bin/python3 /usr/bin/midscroll
# Holds the state socket the session overlay (midscroll-overlay) listens on
RuntimeDirectory=midscroll
Restart=always
RestartSec=2
# Mouse events pass through this process; make sure it's never starved.
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=20

# Sandboxing. The daemon needs only /dev/input, /dev/uinput, the config in
# /etc and its runtime dir. ProtectSystem=strict leaves /dev, /proc and
# /sys writable, and RuntimeDirectory= is implicitly writable, so the
# socket and devices keep working. Deliberately absent because they would
# break the daemon: PrivateDevices/DevicePolicy (needs real evdev+uinput
# nodes, including ones hotplugged later), ProtectKernelModules (the
# ExecStartPre modprobe), RestrictRealtime (FIFO scheduling above).
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
RestrictAddressFamilies=AF_UNIX
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
# Allow only the service syscall set and fail the rest with EPERM rather
# than SIGSYS. @system-service already excludes @resources, @privileged,
# @module, @mount and friends, and still includes ioctl (EVIOCGRAB + uinput)
# and epoll, so grabbing and asyncio keep working. FIFO scheduling is set by
# systemd before exec, so no realtime syscall is needed here.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Pure CPython, no JIT, so it never needs writable-executable memory.
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
